# kx-db-service.slice
[Unit]
Description=KDB-X DB Service Community Edition

[Slice]
MemoryMax=16G
